# systemd/bascula-web.service.d/10-writable-home.conf
[Service]
ProtectHome=read-only
ReadWritePaths=/home/pi/.config/bascula
ReadWritePaths=/run
ReadWritePaths=/tmp
ReadWritePaths=/var/tmp
