# systemd/bascula-web.service
[Unit]
Description=Bascula Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/bascula-cam
EnvironmentFile=-/etc/default/bascula-web
Environment=VENV=/opt/bascula/current/.venv
Environment=APP=/home/pi/bascula-cam
ExecStart=/usr/local/bin/bascula-web
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
