# systemd/bascula-recovery.target
[Unit]
Description=Bascula Recovery Target
Wants=bascula-recovery.service

[Install]
WantedBy=multi-user.target
