# systemd/bascula-recovery.target
[Unit]
Description=Bascula Digital Pro - Recovery Target
Requires=bascula-recovery.service
AllowIsolate=yes
