# systemd/bascula.service
[Unit]
Description=Bascula Digital (Xinit autostart)
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bascula-cam
Environment=PYTHONPATH=/home/pi/bascula-cam
Environment=DISPLAY=:0
Environment=SERIAL_PORT=/dev/serial0
Environment=SERIAL_BAUD=115200
ExecStart=/usr/bin/startx /usr/local/bin/bascula-xsession -- -nocursor
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
