# systemd/bascula-web.service
[Unit]
Description=Bascula Mini-Web (WiFi config)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
EnvironmentFile=-/etc/default/bascula
# Usa el Python del venv:
ExecStart=/opt/bascula/current/.venv/bin/python -m bascula.services.wifi_config --port ${BASCULA_MINIWEB_PORT:-${BASCULA_WEB_PORT:-8080}}
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
