# systemd/bascula-net-fallback.service
[Unit]
Description=Bascula AP fallback (only when no LAN/Wi-Fi connected)
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=oneshot
WorkingDirectory=/opt/bascula/current
EnvironmentFile=-/etc/default/bascula
Environment=VENV=/opt/bascula/current/.venv
Environment=APP=/opt/bascula/current
ExecStart=/bin/bash -lc '/opt/bascula/current/scripts/net-fallback.sh'
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
