# systemd/bascula-app-failure@.service
[Unit]
Description=Registrar fallo consecutivo Báscula App (%i)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/bascula/current
EnvironmentFile=-/etc/default/bascula
Environment=VENV=/opt/bascula/current/.venv
Environment=APP=/opt/bascula/current
Environment=FAIL_UNIT=%i
ExecStart=/bin/bash -lc '/opt/bascula/current/scripts/record_app_failure.sh "${FAIL_UNIT}"'
StandardOutput=journal
StandardError=journal
