# systemd/bascula-ui.service
[Unit]
Description=Bascula UI (Tkinter)
After=network-online.target sound.target

[Service]
User=pi
Group=audio
WorkingDirectory=/home/pi/bascula-cam
# Opcional: forzar tarjeta (comentado por defecto)
# Environment=BASCULA_APLAY_DEVICE=plughw:1,0
Environment=BASCULA_CFG_DIR=%h/.config/bascula
ExecStart=/bin/bash -lc 'if [ -x %h/bascula-cam/.venv/bin/python ]; then %h/bascula-cam/.venv/bin/python %h/bascula-cam/main.py; else python3 %h/bascula-cam/main.py; fi'
Restart=on-failure
RestartSec=3
# Permitir acceso a /dev/snd/*
DeviceAllow=char-116:*

[Install]
WantedBy=graphical.target multi-user.target
