# systemd/bascula-cam.service
# ~/.config/systemd/user/bascula-cam.service
[Unit]
Description=SMART BÁSCULA CAM (Tkinter)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
# Adjust user paths with $HOME expansion
Environment=APP_DIR=%h/bascula-cam
Environment=VENV_DIR=%h/bascula-cam/venv
Environment=DISPLAY=:0
# If needed on your system, also set XAUTHORITY:
# Environment=XAUTHORITY=%h/.Xauthority
WorkingDirectory=%h/bascula-cam
ExecStart=%h/bascula-cam-systemd/scripts/run-bascula.sh
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
