# systemd/bascula-miniweb.service
[Unit]
Description=Bascula Mini-Web (FastAPI)
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=pi
Group=pi
PermissionsStartOnly=true
WorkingDirectory=/opt/bascula/current
Environment=PYTHONPATH=/opt/bascula/current
Environment=BASCULA_ENV=prod
Environment=BASCULA_MINIWEB_PORT=8080
Environment=BASCULA_MINIWEB_STATE=/var/lib/bascula/miniweb
ExecStartPre=/usr/bin/install -d -m 0750 -o pi -g pi /var/lib/bascula/miniweb
ExecStartPre=/usr/bin/install -D -m 0644 -o pi -g pi /dev/null /var/log/bascula/miniweb.log
ExecStart=/opt/bascula/current/.venv/bin/uvicorn bascula.miniweb:app --host 0.0.0.0 --port 8080 --proxy-headers --no-access-log --log-level info
Restart=always
RestartSec=2
StandardOutput=append:/var/log/bascula/miniweb.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
