# systemd/bascula-app.service
[Unit]
Description=Bascula Digital Pro - UI (Xorg via startx)
After=network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service
Conflicts=bascula-recovery.service
OnFailure=
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/bascula/current
EnvironmentFile=-/etc/default/bascula
Environment=HOME=/home/pi
RuntimeDirectory=bascula
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=VENV=/opt/bascula/current/.venv
Environment=APP=/opt/bascula/current
# Ensure log directory and files exist before launching the UI (option A).
ExecStartPre=/bin/bash -lc 'install -d -m 0755 -o pi -g pi /var/log/bascula'
ExecStartPre=/bin/bash -lc 'touch /var/log/bascula/app.log /var/log/bascula/xorg.log'
ExecStartPre=/bin/bash -lc 'chown pi:pi /var/log/bascula/app.log /var/log/bascula/xorg.log'
ExecStartPre=/bin/bash -lc 'chmod 0644 /var/log/bascula/app.log /var/log/bascula/xorg.log'
ExecStartPre=/bin/bash -lc 'test -f /boot/bascula-recovery && { echo "Flag /boot/bascula-recovery detectada" >&2; exit 1; } || true'
ExecStartPre=/bin/bash -lc 'test -f /opt/bascula/shared/userdata/force_recovery && { echo "Flag force_recovery detectada" >&2; exit 1; } || true'
ExecStartPre=/usr/bin/install -d -m 0700 -o pi -g pi /home/pi/.local/share/xorg
ExecStart=/usr/local/bin/bascula-app

Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
