# systemd/x735-poweroff.service
[Unit]
Description=Geekworm x735 v3 safe poweroff monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bascula/current
EnvironmentFile=-/etc/default/bascula
Environment=VENV=/opt/bascula/current/.venv
Environment=APP=/opt/bascula/current
EnvironmentFile=-/etc/default/x735-poweroff
ExecStart=/usr/local/sbin/x735-poweroff.sh
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
