# systemd/x735-fan.service
[Unit]
Description=Geekworm X735 v3 fan controller (libgpiod)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/x735-fan
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
