# system/wifi/wifi-web.service
[Unit]
Description=Bascula WiFi Config Miniweb
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wifi-ap
ExecStart=/usr/bin/python3 -m uvicorn wifi_ap:app --host 0.0.0.0 --port 80
Restart=on-failure

[Install]
WantedBy=multi-user.target
