# packaging/postinst
#!/bin/sh

set -e

LEGACY_DROPIN="/etc/systemd/system/bascula-app.service.d/tty.conf"

if [ -f "$LEGACY_DROPIN" ]; then
    rm -f "$LEGACY_DROPIN"
    rmdir "$(dirname "$LEGACY_DROPIN")" 2>/dev/null || true
    if command -v systemctl >/dev/null 2>&1; then
        systemctl daemon-reload || true
    fi
fi

if command -v loginctl >/dev/null 2>&1 && id pi >/dev/null 2>&1; then
    loginctl enable-linger pi || true
fi

if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload
    systemctl enable --now bascula-app.service x735-fan.service x735-poweroff.service
fi

exit 0
