# etc/systemd/system/bascula-ui.service
[Unit]
Description=Bascula UI (Tkinter Kiosk)
After=graphical.target
Wants=graphical.target

[Service]
User=pi
WorkingDirectory=/home/pi/bascula-cam
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/bascula-cam/scripts/safe_run.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
