# etc/systemd/system/bascula-recovery.service
[Unit]
Description=Bascula UI (Tkinter Kiosk)
After=graphical.target
Wants=graphical.target

[Service]
User=pi
WorkingDirectory=/home/pi/bascula-cam
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc '/home/pi/bascula-cam/.venv/bin/python /home/pi/bascula-cam/main.py --recovery'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
