# etc/systemd/system/bascula-ui.service
[Unit]
Description=Bascula UI (Tkinter Kiosk)
After=graphical.target
Wants=graphical.target

[Service]
User=pi
WorkingDirectory=/home/pi/bascula-cam
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do [ -S /tmp/.X11-unix/X0 ] && exit 0; sleep 0.5; done; exit 1'
ExecStart=/home/pi/bascula-cam/scripts/safe_run.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
