# etc/systemd/system/bascula-miniweb.service
[Unit]
Description=Bascula Mini Web
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/bascula-cam
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/bascula-cam
ExecStartPre=/bin/sh -c '[ -x /home/pi/bascula-cam/.venv/bin/python ]'
ExecStartPre=/bin/sh -c '/home/pi/bascula-cam/.venv/bin/python -c "import uvicorn"'
ExecStart=/home/pi/bascula-cam/scripts/run-miniweb.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
